# student-evaluation-backend.service
[Unit]
Description=Student Evaluation System Backend
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=UCLIXN
WorkingDirectory=/home/UCLIXN/student-evaluation-system/backend
Environment=DATABASE_URL=postgresql://postgres:hdtu4e8JzGqi4I9OUgapRE8Ud@localhost:5432/student_evaluation_system
Environment=REDIS_URL=redis://localhost:6379/0
Environment=CELERY_BROKER_URL=redis://localhost:6379/0
Environment=CELERY_RESULT_BACKEND=redis://localhost:6379/0
EnvironmentFile=/home/UCLIXN/student-evaluation-system/.env
ExecStart=/usr/bin/python3 -m uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target